# crontab
*/10 * * * * /path/to/your/script/DNSCheck.sh